# config/conf.d/haproxy.conf
[program:haproxy]
command=haproxy -f /tmp/haproxy.cfg -p /var/run/haproxy.pid
redirect_stderr=true
numprocs=1
autostart=true
autorestart=unexpected
