# config/conf.d/npm.conf
[program:npm]
command=npm start
redirect_stderr=true
